# 9router.service
[Unit]
Description=9Router AI Proxy
After=network.target

[Service]
ExecStart=/usr/bin/node /usr/lib/node_modules/9router/app/custom-server.js
WorkingDirectory=/root
Environment=NODE_OPTIONS=--max-old-space-size=512
Environment=PORT=20128
Environment=HOSTNAME=0.0.0.0
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
